# supervisord.conf
; supervisord configuration for the factory daemon.
;
; supervisord is the recommended way to run factory when you want it to be able
; to RESTART ITSELF after a self-update: `selfimprove.py` (and the daemon's
; self-update path) call `supervisorctl restart factory`, which only works if a
; supervisord owns the process. (`factory.service` remains a systemd alternative
; for plain always-on running without self-restart.)
;
; Install supervisor and run from this directory:
;   pip install --user supervisor          # or: pipx install supervisor
;   cp .env.example .env && $EDITOR .env    # configure factory first
;   supervisord -c supervisord.conf         # starts the daemon in the background
;   supervisorctl -c supervisord.conf status
;   supervisorctl -c supervisord.conf tail -f factory
;   supervisorctl -c supervisord.conf restart factory
;
; All paths are relative to this file via %(here)s, so the config is portable as
; long as you launch supervisord from the factory checkout.

[unix_http_server]
file=%(here)s/.supervisor/supervisor.sock
chmod=0700

[supervisord]
logfile=%(here)s/.supervisor/supervisord.log
pidfile=%(here)s/.supervisor/supervisord.pid
childlogdir=%(here)s/.supervisor
nodaemon=false
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/.supervisor/supervisor.sock

[program:factory]
command=%(here)s/run.sh
directory=%(here)s
autostart=true
autorestart=true
startsecs=5
startretries=3
; factory finishes the current work item on SIGTERM before exiting, so give the
; graceful stop (and any in-flight `printer exec`) plenty of room.
stopsignal=TERM
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=%(here)s/.supervisor/factory.out.log
stderr_logfile=%(here)s/.supervisor/factory.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
; factory drives `claude`, `printer`, `git`, `gh`, and `supervisorctl`; keep a
; login-like PATH so those binaries resolve under the daemon. Adjust to match
; your machine (mirrors the PATH in factory.service).
environment=PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin"
